# dbus/bus/dbus.service.in
[Unit]
Description=D-Bus System Message Bus
Requires=dbus.socket
DefaultDependencies=no
After=local-fs.target dbus.socket
Before=basic.target
Conflicts=shutdown.target

[Service]
ExecStartPre=@EXPANDED_BINDIR@/dbus-uuidgen --ensure
ExecStart=@EXPANDED_BINDIR@/dbus-daemon --system --address=systemd: --nofork --nopidfile --systemd-activation
ExecReload=@EXPANDED_BINDIR@/dbus-send --print-reply --system --type=method_call --dest=org.freedesktop.DBus / org.freedesktop.DBus.ReloadConfig
OOMScoreAdjust=-900

[Install]
WantedBy=basic.target
